# scripts/trainpixels.service
[Unit]
Description=TrainPixels LED Controller Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=trainpixels
Group=trainpixels
WorkingDirectory=/opt/trainpixels
ExecStart=/opt/trainpixels/trainpixels
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/bin:/bin:/usr/local/bin
Environment=PYTHONPATH=/opt/trainpixels

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/trainpixels
ReadWritePaths=/home/trainpixels/.config/trainpixels

[Install]
WantedBy=multi-user.target